# configs/openssh/supervisord-openssh.conf
[program:openssh]
command=/start-openssh.sh
