# configs/openssh/supervisord-openssh.conf
[supervisord]
nodaemon=true

[program:openssh]
command=/start-openssh.sh
