# configs/mysql/supervisord-mysqld.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/start-mysqld.sh
numprocs=4
process_name=%(program_name)s_%(process_num)
autostart=true
autorestart=true
